# rules.d/90-security.rules
ACTION=="add",KERNEL=="sdb*",RUN+="/code/security_udev/udev.sh"
